# supervisor/conf.d/chinadns.ini
[program:chinadns]
command=chinadns -c /root/faure.sh/data/chnlist.txt -s 223.5.5.5,9.9.9.9,8.8.4.4 -p 5353 -b 127.0.0.1
user=root
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/china_dns.log
loglevel=info
