# supervisor/conf.d/chinadns.ini
[program:chinadns]
command=chinadns -c /root/faure.sh/data/chnlist.txt -p 5353 -b 127.0.0.1
user=root
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/china_dns.log
loglevel=info
