# supervisor/conf.d/http-proxy.ini
[program:http-proxy]
command=gost -L=http://:1081 -F socks5://127.0.0.1:1080
user=nobody
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/http-proxy.log
loglevel=info
