# supervisor/conf.d/redirect.ini
[program:redirect]
command=gost -L redirect://:8848 -F socks5://127.0.0.1:1080
user=root
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/redirect.log
loglevel=info
