# supervisor/conf.d/ss-local.ini
[program:ss-local]
command=ss-local -c /root/shadowsocksr.json -u --fast-open --mptcp -t 30 -l 1080 -b 0.0.0.0
user=root
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/ss-local.log
loglevel=info
